# enso-kiosk.service
[Unit]
Description=Enso Karate Catalog Kiosk Mode
After=network.target enso-catalog.service
Requires=enso-catalog.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
ExecStartPre=/home/pi/enso-catalog/map-touchscreen.sh
ExecStart=/home/pi/enso-catalog/kiosk-launcher.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target