# enso-kiosk.service
[Unit]
Description=Enso Karate Catalog Kiosk Mode
After=network.target enso-catalog.service
Requires=enso-catalog.service

[Service]
Type=simple
# Run as root and use wrapper to determine session user at runtime
# (wrapper will sudo to the detected user)
ExecStart=/home/pi/enso-catalog/kiosk-start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target