# enso-kiosk.service
[Unit]
Description=Enso Karate Catalog Kiosk Mode
After=network.target enso-catalog.service
Requires=enso-catalog.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
ExecStart=xdg-open http://localhost:5000/kiosk
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target