# enso-kiosk.service
[Unit]
Description=Enso Karate Catalog Kiosk Mode
After=network.target

[Service]
Type=simple
ExecStart=@INSTALL_DIR@/kiosk-launcher.sh
Restart=no

[Install]
WantedBy=default.target