# enso-kiosk.service
[Unit]
Description=Enso Karate Catalog Kiosk Mode
After=network.target enso-catalog.service
Requires=enso-catalog.service

[Service]
Type=simple
User=enso
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/enso/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/home/enso/enso-catalog/kiosk-launcher.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target